# melon-stream-proxy.service
[Unit]
Description=Melon Stream Proxy
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/melon-stream-proxy
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
